# config/rc.local
#! /bin/bash

if [ -e /etc/configboot ]; then
  hostnamectl set-hostname rock4se
  chmod u+s /usr/bin/sudo
  if dpkg-query -W "lightdm" >/dev/null 2>&1; then
    dpkg-reconfigure lightdm
  fi
  if dpkg-query -W "sddm" >/dev/null 2>&1; then
    dpkg-reconfigure sddm
  fi
  if dpkg-query -W "gdm3" >/dev/null 2>&1; then
    dpkg-reconfigure gdm3
  fi
  rm /etc/configboot
  touch /etc/default/FIRSTBOOT
  reboot
fi

if [ -e /etc/default/FIRSTBOOT ]; then
  resizeroot
fi


